# dump/dist/pvm_dump_offload@.service.in
[Unit]
Description=PowerVM Handler Host %i
Wants=obmc-host-started@%i.target
After=xyz.openbmc_project.Dump.Manager.service
After=org.open_power.Dump.Manager.service
After=xyz.openbmc_project.biosconfig_manager.service
After=obmc-host-started@%i.target
Conflicts=obmc-host-stop@%i.target

[Service]
ExecStart=@bindir@/pvm_dump_offload
Restart=on-failure
SyslogIdentifier=pvm_dump_offload

[Install]
WantedBy=obmc-host-startmin@0.target
