# dump/dist/pvm_dump_offload.service.in
[Unit]
Description=PowerVM Handler
Wants=phosphor-debug-collector.service
After=phosphor-debug-collector.service
After=xyz.openbmc_project.biosconfig_manager.service
After=obmc-power-on@0.target
Conflicts=obmc-host-stop@0.target
Conflicts=obmc-host-timeout@0.target

[Service]
ExecStart=@bindir@/pvm_dump_offload
Restart=on-failure
SyslogIdentifier=pvm_dump_offload

[Install]
WantedBy=obmc-host-startmin@0.target
